# os/linux/aide/aidecheck.timer
[Unit]
Description=Daily AIDE check timer

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target