# os/linux/system-monitoring/archive/linux-sentinel.service
[Unit]
Description=Linux Sentinel v2.3 Security Monitor
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/linux-sentinel/linux_sentinel.py enhanced
User=root
StandardOutput=journal
StandardError=journal
Restart=on-failure
CapabilityBoundingSet=CAP_NET_ADMIN CAP_DAC_READ_SEARCH CAP_SYS_PTRACE
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target