# os/linux/system-monitoring/rust/linux-sentinel.timer
[Unit]
Description=Run Linux Sentinel hourly
Requires=linux-sentinel.service

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target