# detect_rules/splunk/import/searches/interlock_ransomware_activity.json
[
  {
    "name": "Interlock Malware Detection",
    "search": "| tstats `summariesonly` count min(_time) as firstTime max(_time) as lastTime from datamodel=Endpoint.Processes where `interlock_hashes` by Processes.process_name Processes.process Processes.parent_process Processes.dest Processes.user Processes.sha256 | `drop_dm_object_name(\"Processes\")` | eval Tactic=\"Execution\", Technique=\"T1204.002\", DetectionMethod=\"Known Malicious Hash\" | append [ | tstats `summariesonly` count min(_time) as firstTime max(_time) as lastTime from datamodel=Endpoint.Processes where (Processes.process_name=powershell.exe OR Processes.process=\"*powershell*\") AND (Processes.process=\"*irm *\" OR Processes.process=\"*iex *\" OR Processes.process=\"*Invoke-RestMethod*\" OR Processes.process=\"*Invoke-Expression*\" OR Processes.process=\"*-w h*\" OR Processes.process=\"*-windowstyle hidden*\") by Processes.process_name Processes.process Processes.parent_process Processes.dest Processes.user Processes.sha256 | `drop_dm_object_name(\"Processes\")` | eval Tactic=\"Execution\", Technique=\"T1059.001\", DetectionMethod=\"Suspicious PowerShell Command\" ] | append [ | tstats `summariesonly` count min(_time) as firstTime max(_time) as lastTime from datamodel=Endpoint.Registry where Registry.registry_path=\"*\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run*\" AND `interlock_reg_values` by Registry.registry_path Registry.registry_value_name Registry.registry_value_data Registry.process_guid Registry.dest Registry.user | `drop_dm_object_name(\"Registry\")` | rename process_guid as guid | join type=left guid [ | tstats `summariesonly` count from datamodel=Endpoint.Processes where earliest=-1d by Processes.process_name Processes.process Processes.parent_process Processes.guid | `drop_dm_object_name(\"Processes\")` ] | eval Tactic=\"Persistence\", Technique=\"T1547.001\", DetectionMethod=\"Registry Run Key Modification\" ] | append [ | tstats `summariesonly` count min(_time) as firstTime max(_time) as lastTime from datamodel=Network_Traffic.All_Traffic where (`interlock_c2_ips` OR `interlock_domains` OR dest_host=\"*trycloudflare.com*\") by All_Traffic.src_ip All_Traffic.dest_ip All_Traffic.dest_host All_Traffic.process_guid All_Traffic.dest All_Traffic.user | `drop_dm_object_name(\"All_Traffic\")` | rename process_guid as guid | join type=left guid [ | tstats `summariesonly` count from datamodel=Endpoint.Processes where earliest=-1d by Processes.process_name Processes.process Processes.parent_process Processes.guid | `drop_dm_object_name(\"Processes\")` ] | eval Tactic=\"Command and Control\", Technique=\"T1071.001\", DetectionMethod=\"C2 Communication\" ] | append [ | tstats `summariesonly` count min(_time) as firstTime max(_time) as lastTime from datamodel=Endpoint.Processes where Processes.process_name=schtasks.exe AND Processes.process=\"*/create*\" AND (Processes.process=\"*cmd*\" OR Processes.process=\"*powershell*\") by Processes.process_name Processes.process Processes.parent_process Processes.dest Processes.user Processes.sha256 | `drop_dm_object_name(\"Processes\")` | eval Tactic=\"Persistence\", Technique=\"T1053.005\", DetectionMethod=\"Scheduled Task Creation\" ] | rename dest as DeviceName, process_name as FileName, process as ProcessCommandLine, parent_process as InitiatingProcess, sha256 as SHA256, registry_path as RegistryKey, registry_value_name as RegistryValueName, registry_value_data as RegistryValueData, src_ip as SourceIP, dest_ip as DestinationIP, dest_host as DestinationHost | table firstTime, lastTime, DeviceName, user, Tactic, Technique, DetectionMethod, FileName, ProcessCommandLine, InitiatingProcess, SHA256, RegistryKey, RegistryValueName, RegistryValueData, SourceIP, DestinationIP, DestinationHost",
    "description": "Correlates various indicators and techniques associated with Interlock malware, including malicious hashes, PowerShell execution, registry run key persistence, C2 communication, and scheduled task creation.",
    "is_scheduled": 1,
    "disabled": 0,
    "cron_schedule": "*/30 * * * *",
    "dispatch.earliest_time": "-24h",
    "dispatch.latest_time": "now",
    "sharing": "global",
    "acl": {
      "read": ["power", "admin"],
      "write": ["admin"]
    }
  }
]
