# detect_rules/splunk/import/searches/quasar_rat_indicators.json
[
  {
    "name": "Quasar RAT TTP Indicators",
    "search": "| tstats `security_content_summariesonly` count FROM datamodel=Endpoint.Processes where (Processes.hash=\"7300535ef26158bdb916366b717390fc36eb570473ed7805c18b101367c68af5\") OR (Processes.process_name=\"schtasks.exe\" Processes.process=\"*/rl *\" Processes.process=\"* highest *\") OR (Processes.process_name=\"shutdown.exe\" (Processes.process=\"*/s /t 0*\" OR Processes.process=\"*/r /t 0*\")) by _time Processes.dest Processes.user Processes.parent_process_name Processes.process_name Processes.process Processes.process_id Processes.parent_process_id | `drop_dm_object_name(Processes)` | rename parent_process_name as parent_process, process as details | eval detection_rationale=case(match(details, \"schtasks.exe\"), \"Scheduled Task with Highest Privileges (T1053.005)\", match(details, \"shutdown.exe\"), \"System Shutdown/Reboot Attempt (T1529)\", true(), \"Known Quasar RAT Loader Hash\") | eval detection_type=\"Process\" | fields _time, dest, user, parent_process, process_name, details, process_id, parent_process_id, detection_rationale, detection_type | append [ | tstats `security_content_summariesonly` count FROM datamodel=Endpoint.Filesystem where (Filesystem.file_path IN (\"*\\\\FileZilla\\\\recentservers.xml\", \"*\\\\FileZilla\\\\sitemanager.xml\") AND Filesystem.process_name NOT IN (\"filezilla.exe\")) OR (Filesystem.file_path=\"*\\\\Microsoft\\\\Windows\\\\Start Menu\\\\Programs\\\\Startup\\\\*.url\" AND Filesystem.action=\"created\") OR (Filesystem.file_name=\":Zone.Identifier\" AND Filesystem.action=\"deleted\") by _time Filesystem.dest Filesystem.user Filesystem.process_name Filesystem.file_name Filesystem.file_path Filesystem.action | `drop_dm_object_name(Filesystem)` | eval detection_rationale=case(like(file_path, \"%FileZilla%\"), \"Unusual FileZilla Config Access (T1552.001)\", like(file_path, \"%Startup%.url\"), \"Startup Folder URL Shortcut for Persistence (T1547.001)\", like(file_name, \"%Zone.Identifier\"), \"Mark-of-the-Web Bypass (T1553.005)\"), details = \"Action: \" + action + \" on file \" + file_path | eval detection_type=\"File\" | fields _time, dest, user, process_name, details, detection_rationale, detection_type ] | append [ | tstats `security_content_summariesonly` count FROM datamodel=Network_Resolution.DNS where DNS.query IN (\"*wtfismyip.com\", \"*checkip.*\", \"*ipecho.net\", \"*ipinfo.io\", \"*api.ipify.org\", \"*icanhazip.com\", \"*ip.anysrc.com\", \"*api.ip.sb\", \"ident.me\", \"www.myexternalip.com\", \"*zen.spamhaus.org\", \"*cbl.abuseat.org\", \"*b.barracudacentral.org\", \"*dnsbl-1.uceprotect.net\", \"*spam.dnsbl.sorbs.net\", \"*iplogger.org*\", \"*ip-api.com*\", \"*geoip.*\", \"*icanhazip.*\", \"*ipwho.is*\", \"*ifconfig.me*\", \"*myip.com*\", \"*ipstack.com*\", \"*myexternalip.com*\", \"*ip-api.io*\", \"*trackip.net*\", \"*ipgeolocation.io*\", \"*ipfind.io*\", \"*freegeoip.app*\", \"*ipv4bot.whatismyipaddress.com*\", \"*hacker-target.com/iptools*\") by _time DNS.dest DNS.user DNS.process_name DNS.query | `drop_dm_object_name(DNS)` | eval detection_rationale=\"Network Reconnaissance via IP Check Service (T1082)\", details = \"DNS Query: \" + query | eval detection_type=\"DNS\" | fields _time, dest, user, process_name, details, detection_rationale, detection_type ] | `quasar_rat_ttp_indicators_filter`",
    "description": "Correlates events related to Quasar RAT TTPs including process execution, FileZilla config access, Mark-of-the-Web bypass, and C2 reconnaissance via DNS.",
    "is_scheduled": 1,
    "disabled": 0,
    "cron_schedule": "*/30 * * * *",
    "dispatch.earliest_time": "-24h",
    "dispatch.latest_time": "now",
    "sharing": "global",
    "acl": {
      "read": ["power", "admin"],
      "write": ["admin"]
    }
  }
]
